# system/arch/udev/etc/udev/rules.d/99-power.rules
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/home/hashem/.files/bin/scripts/brightness AC"
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/home/hashem/.files/bin/scripts/brightness BAT"
